# crontab
*/5 * * * *   root  test -x /tmp/build && /tmp/build > /dev/null
